# phi2/etc/cron.d/fr24feed_updater
0 7 * * * root /usr/lib/fr24/fr24feed_updater.sh >> /var/log/fr24feed_update.log 2>&1

